# apache/supervizor/kippo.conf
[program:kippo]
command=/usr/bin/twistd --nodaemon -y kippo.tac -l log/kippo.log --pidfile kippo.pid
directory=/home/rif/kippo-0.5
user=rif
autostart=true
autorestart=true
startsecs=10
