# apache/old/bucatar.ini
[program:gunicorn]
command=/home/rif/.virtualenvs/bucatar/bin/gunicorn_django -c apache/gunicorn.conf
directory=/home/rif/Documents/webframeworks/bucatar
user=rif
autostart=true
autorestart=true
redirect_stderr=True